# deployment/chronochunk.service
[Unit]
Description=ChronoChunk Discord Bot
After=network.target

[Service]
Type=simple
User=mykal
WorkingDirectory=/opt/chronochunk
ExecStart=/opt/chronochunk/venv/bin/python run.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=chronochunk
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target